# Debian/schedule/daily_iPhoneBackup.service
[Unit]
Description = "Daily iPhone Backup"

[Service]
ExecStart = /home/megabits/scripts/daily_iPhoneBackup.sh
