# Debian/schedule/daily_archive.service
[Unit]
Description=Daily Archive
[Service]
ExecStart=/home/megabits/scripts/daily_archive.sh
