# Debian/schedule/daily_iPhoneBackup.timer
[Unit]
Description=Run every day

[Timer]
OnCalendar=*-*-* 23:00:00

[Install]
WantedBy=timers.target
