# Debian/schedule/monthly_clearArchive.service
[Unit]
Description = "Monthly Clear Archive"

[Service]
ExecStart = /home/megabits/scripts/monthly_clearArchive.sh
