# Debian/schedule/monthly_clearArchive.timer
[Unit]
Description=Run every month

[Timer]
OnCalendar=*-*-1 05:00:00

[Install]
WantedBy=timers.target
